# dashcam.service
[Unit]
Description=Dashcam Recorder
After=multi-user.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/dashcam
ExecStart=/usr/bin/python3 /home/pi/dashcam/recorder.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
